# systemd/renew-all-users.timer
[Unit]
Description=Timer for renewing all active users' books

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
